# service/log/run
#!/bin/sh
exec 2>&1
mkdir -p /var/log/dbus-conext-bridge
exec multilog t s25000 n4 /var/log/dbus-conext-bridge
